# RADXAQ6A/image/files/systemd/can0.service
[Unit]
Description=Bring up can0 at 500 kbit/s (Waveshare RS485 CAN HAT, MCP2515 on SPI12/CS0)
# We deliberately do NOT use ConditionPathExists=/sys/class/net/can0 here.
# That condition is evaluated ONCE at service start time, which on a cold
# boot regularly fires before the MCP2515 SPI driver has finished binding
# and creating the netdev — leading to the service being silently skipped
# and the interface staying down despite a perfectly functional HAT.
#
# We also do NOT use Wants=/After=sys-subsystem-net-devices-can0.device:
# that blocks the boot graph ~90 s on systems where no HAT is installed
# (bench / dev / manufacturing).
#
# Instead, headwaters-can0-up.sh polls for /sys/class/net/can0 with a
# bounded 30 s deadline, sets bitrate + brings up the interface when it
# appears, and exits cleanly without error if it doesn't. Race-free,
# never blocks boot.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/headwaters-can0-up.sh
ExecStop=/sbin/ip link set can0 down
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
