# RADXAQ6A/image/files/systemd/cpu-powersave.service
[Unit]
Description=Cap all CPU cores to 600 MHz with powersave governor (CM5-parity target)
After=multi-user.target

# The QCS6490 has 4× Kryo Gold Prime (A78) cores @ 2.7 GHz and 4× Kryo Silver
# (A55) cores @ 1.9 GHz — 8 cores total.
#
# The CM5 image caps ARM at `arm_freq=600` and ships the same Headwaters
# Docker workload (Mosquitto + Mongo + Node + nginx + tileserver) with no
# trouble. We match that ceiling exactly on every Q6A core, which:
#
#   * Gives the A55 little cores the same per-core headroom as the CM5's
#     A76 cores at 600 MHz (the A55 at 600 MHz is weaker per cycle, but we
#     have 4 of them plus the 4 big cores, so total compute is ≥ CM5).
#
#   * Leaves all 4 A78 big cores ONLINE so bursty work (Docker image pulls,
#     tile rendering, OTA unpack) can migrate to them and finish at least
#     as fast as on the CM5. We explicitly do NOT offline the big cores —
#     that would drop burst capacity below CM5 parity.
#
#   * Keeps the A78 cores pinned to 600 MHz max, so when idle they sit in
#     their deepest C-state most of the time and contribute minimally to
#     idle power. The powersave governor also prevents any opportunistic
#     boosting.
#
# The picked frequency is the highest available step ≤ 600 MHz on each
# policy group. If 600 MHz is not a supported step, the kernel rounds
# down to the nearest lower one.

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  TARGET_KHZ=600000; \
  for cpu in /sys/devices/system/cpu/cpu[0-9]*; do \
    [ -d "$cpu/cpufreq" ] || continue; \
    echo powersave > "$cpu/cpufreq/scaling_governor" 2>/dev/null || true; \
    avail=$(cat "$cpu/cpufreq/scaling_available_frequencies" 2>/dev/null); \
    pick=""; \
    for f in $avail; do \
      if [ "$f" -le "$TARGET_KHZ" ]; then \
        [ -z "$pick" ] && pick="$f"; \
        [ "$f" -gt "$pick" ] && pick="$f"; \
      fi; \
    done; \
    [ -z "$pick" ] && pick="$TARGET_KHZ"; \
    echo "$pick" > "$cpu/cpufreq/scaling_max_freq" 2>/dev/null || true; \
  done; \
  true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
