# RADXAQ6A/image/files/systemd/headwaters-firstboot-network.service
[Unit]
Description=TrailCurrent Headwaters first-boot setup (network stage: Python venv + pip)
# Run only after the early firstboot finished its non-network work.
Wants=headwaters-firstboot.service
After=headwaters-firstboot.service
# We need DNS + HTTPS to PyPI. network-online.target is the canonical gate.
Wants=network-online.target
After=network-online.target
# Don't block sysinit on this — it's a deferred, late-boot job.
DefaultDependencies=yes
# Re-run on every boot until the sentinel exists. The script writes the
# sentinel only after pip install succeeds AND import verification passes,
# so a partial failure means the next boot retries automatically.
ConditionPathExists=!/var/lib/headwaters/.firstboot-network-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/headwaters-firstboot-network.sh
StandardOutput=journal+console
StandardError=journal+console
# Generous timeout: pip on a fresh system pulling ~8 wheels for arm64
# from PyPI typically completes in 60-90s, but a slow first-flight TLS
# handshake or wheel cache miss can push it to several minutes. The
# script's own retry loop already backs off; this is the outer cap.
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
