# RADXAQ6A/image/files/systemd/headwaters-load-images.service
[Unit]
Description=Load baked-in Docker images from tarballs (one-shot)
After=docker.service headwaters-firstboot.service
Requires=docker.service
ConditionPathExists=/home/trailcurrent/images
ConditionDirectoryNotEmpty=/home/trailcurrent/images

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/headwaters-load-images.sh
RemainAfterExit=yes
TimeoutStartSec=600
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
