# RADXAQ6A/image/files/systemd/headwaters-firstboot.service
[Unit]
Description=TrailCurrent Headwaters first-boot setup (expand NVMe, regen IDs, TLS certs, Python venv)
DefaultDependencies=no
After=systemd-remount-fs.service local-fs.target
Before=sysinit.target shutdown.target docker.service
Conflicts=shutdown.target
ConditionPathExists=!/var/lib/headwaters/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/headwaters-firstboot.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=sysinit.target
